# 2020.reorder/local/rc.local
#!/bin/sh

# mtu black hole
iptables -I FORWARD -p tcp --tcp-flags SYN,RST SYN -j TCPMSS --clamp-mss-to-pmtu

