# 2023/abandoned/elementary.io/loki/add-on/etc/rc.local
#!/bin/sh -e
#
# rc.local
#

powertop --auto-tune

exit 0
