# 2023/abandoned/elementary.io/loki/add-on/home/.config/autostart/custom-command1.desktop
[Desktop Entry]
Name[en_US]=Custom Command
Comment[en_US]=/usr/bin/rescuetime 
Exec=/usr/bin/rescuetime 
Icon=application-default-icon
X-GNOME-Autostart-enabled=true
Type=Application
