# 2020.reorder/network/fwtraf/fwtraf.cron
# Минута  Час     День месяца  Месяц   День недели  Пользователь             Задание
# (0-59)  (0-23)  (1-31)       (1-12)  (0-7)
# ------  ------  -----------  ------  -----------  ------------  ---------  --------------------------------------------
  58      *       *            *       *            root                     /bin/fwtraf trafhourly; /bin/fwtraf trafmake
  30      23      *            *       0            root                     /bin/fwtraf ipdbclear
